# deployment/system/easeon.socket
[Unit]
Description=easeon socket

[Socket]
ListenStream=/run/easeon.sock

[Install]
WantedBy=sockets.target
