# deployment/gunicorn/easeOn.socket
[Unit]
Description=easeOn socket

[Socket]
ListenStream=/run/easeOn.sock

[Install]
WantedBy=sockets.target
